# systemd/ser2net.service
[Unit]
Description=ser2net — serial-to-network bridge with web config UI
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Run as a dedicated unprivileged user that is a member of the 'dialout' group
# so it can open serial devices (NEVER run this as root).
User=ser2net
SupplementaryGroups=dialout
WorkingDirectory=/opt/ser2net
# On a headless first boot the console picker is non-interactive and defaults the
# admin UI to 127.0.0.1:8080 (loopback, never network-exposed without an explicit
# choice). To expose it on the LAN, run once interactively (or with --reconfigure)
# to pick the bind IP and set a password, or pre-create /opt/ser2net/data/config.json.
ExecStart=/usr/bin/python3 /opt/ser2net/ser2net.py --data-dir /opt/ser2net/data
Restart=on-failure
RestartSec=3

# --- Hardening ---
NoNewPrivileges=true
# Writable paths (rest of FS is read-only under ProtectSystem=strict):
#  - /opt/ser2net/data : config.json + logs (adjust to your --data-dir)
#  - /opt/ser2net/lib  : where the first-run offline bootstrap installs the bundled
#    wheels. Without it ProtectSystem=strict makes ./lib read-only and the first
#    start fails to install dependencies. If you instead pre-install deps at deploy
#    time (python3 bootstrap.py) and add --no-bootstrap to ExecStart, you can drop
#    /opt/ser2net/lib from this list.
ReadWritePaths=/opt/ser2net/data /opt/ser2net/lib
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectProc=invisible
ProcSubset=pid
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectClock=true
ProtectHostname=true
RestrictNamespaces=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
# Add AF_NETLINK to enable pyudev event-driven hotplug; without it ser2net falls
# back to polling (which works fine), so it is left out of the default allow-list.
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
# MemoryDenyWriteExecute=true   # enable only after verifying ctypes/Python startup
# UMask=0077

[Install]
WantedBy=multi-user.target
